# provision/systemd/nullclaw.service
[Unit]
Description=nullclaw gateway (personal agent, discord channel)
After=network-online.target lightpanda-mcp.service
Wants=network-online.target lightpanda-mcp.service

[Service]
User=tetraslam
# mise shims so mcp stdio children (npx for context7) resolve
Environment=PATH=/home/tetraslam/.local/share/mise/shims:/home/tetraslam/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/nullclaw gateway
WorkingDirectory=/home/tetraslam
Restart=always
RestartSec=5
MemoryMax=2G
# WORKAROUND (upstream #977): some gateway connections are born deaf, so
# recycle at most every 10 min even without traffic. delete with the defib.
RuntimeMaxSec=600

[Install]
WantedBy=multi-user.target
